# systemd/user/first-login.service
[Unit]
Description=Verify the first configured Hyprland session
PartOf=graphical-session.target
After=monitor-setup.service waybar.service swaync.service graphical-session.target
ConditionPathExists=!%h/.local/state/dotfiles/session-verified

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 5
ExecStart=%h/.local/bin/first-login-check

[Install]
WantedBy=graphical-session.target
